# src/docker/conf/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:httpd]
command=/usr/local/bin/httpd-foreground
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
